# tests/fixtures/startup/agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>org.example.safe.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Example.app/Contents/MacOS/Example</string>
    <string>--background</string>
  </array>
  <key>AssociatedBundleIdentifiers</key>
  <array><string>org.example.safe</string></array>
  <key>RunAtLoad</key><true/>
</dict>
</plist>
